# install/etc/services.d/postfix/run
#!/usr/bin/with-contenv bash

/usr/libexec/postfix/master -d